# slurm_inspector.conf
# slurm_inspector
# Ubuntu Upstart service configuration
# Copy this file to /etc/init/ and make changes if necessary.
#
# To check syntax of this script:
# init-checkconf /etc/init/slurm_inspector.conf
#
# To reload this script after changes:
# sudo initctl reload-configuration
#
# To start the slurm_inspector service:
# sudo start slurm_inspector

description "slurm_inspector, monitors SLURM and present results via web page"
author "Willi Kappler"

# Stanzas
#
# Stanzas control when and how a process is started and stopped
# See a list of stanzas here: http://upstart.ubuntu.com/wiki/Stanzas

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn

# Specify working directory, change if needed
chdir /var/www/slurm_inspector

# Specify the process/command to start, add command line options if needed
# (Use "slurm_inspector --help" to get a list of options)
exec ./slurm_inspector
